# ansible/roles/pgwc/templates/pgwc.service.j2
[Unit]
Description= pgwc service

[Service]
ExecStart={{ SERVICE_SCRIPT_PATH }}

[Install]
WantedBy=default.target
