# ansible/roles/sgwc/templates/sgwc.service.j2
[Unit]
Description= sgwc service

[Service]
ExecStart={{ SERVICE_SCRIPT_PATH }}

[Install]
WantedBy=default.target
