# ansible/roles/spgwc/templates/spgwc.service.j2
[Unit]
Description= spgwc service

[Service]
ExecStart={{ SERVICE_SCRIPT_PATH }}

[Install]
WantedBy=default.target
